# base/apache_conf/supervisord.httpd.conf
[program:httpd]
command=apachectl -D FOREGROUND
redirect_stderr=true
killasgroup=true
stopasgroup=true
